# rtc-ds1307.conf
# Ubuntu upstart file at /etc/init/rtc-ds1307.conf
# Enables DS1307 RTC Breakout on Ubuntu 14.04 armhf on RPi2

start on runlevel [2345] stop on runlevel [06]

script
# Redirect output to kernel log
#exec >/dev/kmsg 2>&1 # uncomment to enable debug output
#echo “rtc-ds1307.conf trying to load” # debug output

# Script iterates i2c buses and tries to
# find DS1307 clock. If successful then
# sets linux system time from RTC.

# Inteneded to be used on Ubuntu 14.04 armhf
# On Raspberry Pi 2 B+

# This script should be called from by
# Upstart from /etc/init/rtc-ds1307
# on entering runlevel [2345]

# If clock already exist then do nothing
if [ -e /dev/rtc0 ] then
return 0
fi

# Iterate all i2c buses (more generic)
for bus in $(ls -d /sys/bus/i2c/devices/i2c-*);
do
# Try ds1307 and wait for it to initialise
echo ds1307 0x68 >> $bus/new_device;
sleep 0.5
if [ -e /dev/rtc0 ];
then
# DS1307 found, set sys clock from rtc
#echo “rtc found at ” $bus # debug output
hwclock -s -f /dev/rtc0
break # We are done
else
# RTC not found, unregister i2c device
#echo “no rtc at ” $bus # debug output
echo 0x68 >> $bus/delete_device
fi
done

end script